# deploy/trading-daemon.service
[Unit]
Description=Trading System Monitoring Daemon
Documentation=https://github.com/IgorGanapolsky/trading
After=network.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading
Environment=PYTHONPATH=/opt/trading
Environment=GITHUB_TOKEN=your_github_token_here
Environment=SLACK_WEBHOOK=your_slack_webhook_here

ExecStart=/usr/bin/python3 -m src.monitoring.trading_daemon --interval 300
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=30

# Logging
StandardOutput=append:/var/log/trading/daemon.log
StandardError=append:/var/log/trading/daemon_error.log

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trading/data /opt/trading/logs /var/log/trading

[Install]
WantedBy=multi-user.target
