# launchd/com.trading.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trading.healthcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/igorganapolsky/workspace/git/apps/trading/venv/bin/python</string>
        <string>/Users/igorganapolsky/workspace/git/apps/trading/scripts/health_check.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>5</integer>
            <key>Weekday</key>
            <integer>1</integer>  <!-- Monday -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>5</integer>
            <key>Weekday</key>
            <integer>2</integer>  <!-- Tuesday -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>5</integer>
            <key>Weekday</key>
            <integer>3</integer>  <!-- Wednesday -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>5</integer>
            <key>Weekday</key>
            <integer>4</integer>  <!-- Thursday -->
        </dict>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>5</integer>
            <key>Weekday</key>
            <integer>5</integer>  <!-- Friday -->
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/igorganapolsky/workspace/git/apps/trading</string>

    <key>StandardOutPath</key>
    <string>/Users/igorganapolsky/workspace/git/apps/trading/logs/health_check_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/igorganapolsky/workspace/git/apps/trading/logs/health_check_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
